# support-files/templates/#etc#supervisor-bklog-grafana.conf
[unix_http_server]
file=/var/run/bklog/grafana-supervisor.sock             ; (the path to the socket file)

[supervisord]
logfile=__BK_HOME__/logs/bklog/grafana-supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/bklog/grafana-supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
environment=LOGGER_WITHOUT_CONSOLE=1
directory=__BK_HOME__/bklog/grafana

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/bklog/grafana-supervisor.sock

[program:grafana]
process_name=grafana
command=__BK_HOME__/bklog/grafana/bin/grafana-server --homepath ./ --config ./conf/grafana.ini cfg:default.paths.logs=__BK_HOME__/logs/bklog/ cfg:default.paths.data=./data cfg:default.paths.plugins=./plugins
numprocs=1                    ; number of processes copies to start (def 1)
priority=520                  ; the relative start priority (default 999)
startsecs=0                   ; number of secs prog must stay running (def. 1)
stopwaitsecs=30
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=__BK_HOME__/logs/bklog/grafana-error.log
redirect_stderr=false
stopasgroup=true
killasgroup=true
